# deploy/stockpro-backend.service
[Unit]
Description=StockPro Backend API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/stockpro/backend
EnvironmentFile=/opt/stockpro/backend/.env
ExecStart=/opt/stockpro/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 4445
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
